# scripts/mac-library-backup/com.immich.mac-library-backup.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mac Photos originals → delta NFS (weekly incremental).
  Install:
    IMMICH_APPS_ROOT=/path/to/immich-apps
    sed "s|__IMMICH_APPS_ROOT__|$IMMICH_APPS_ROOT|g; s|__HOME__|$HOME|g" \
      scripts/mac-library-backup/com.immich.mac-library-backup.plist.example \
      > ~/Library/LaunchAgents/com.immich.mac-library-backup.plist
    launchctl load ~/Library/LaunchAgents/com.immich.mac-library-backup.plist

  Schedule: Saturday 02:00 (錯開週日 04:00 immich-data-backup on cluster).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immich.mac-library-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__IMMICH_APPS_ROOT__/scripts/mac-library-backup-rsync.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/immich-mac-backup/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/immich-mac-backup/launchd.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.pyenv/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
