# scripts/photo-sync/com.immich.photo-sync.reconcile-watch.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immich.photo-sync.reconcile-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__IMMICH_APPS_ROOT__/scripts/photo-sync/run-with-op-env.sh</string>
        <string>__IMMICH_APPS_ROOT__/scripts/photo-sync/immich-reconcile-watch.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/immich-photo-sync/reconcile.watch.launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/immich-photo-sync/reconcile.watch.launchd.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.pyenv/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PHOTO_SYNC_CONFIG</key>
        <string>__HOME__/.config/immich-apps/photo-sync.config.yaml</string>
    </dict>
</dict>
</plist>
